# services/supervisor/conf.d/php56-fpm.ini
[program:php-fpm56]
command=/usr/local/sbin/php-fpm

priority=99                ; the relative start priority (default 999)
autostart=true              ; start at supervisord start (default: true)
autorestart=true            ; retstart at unexpected quit (default: true)
user=root                 ; setuid to this UNIX account to run the program

stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0